# configs/kvm/usbhotplug/udev/90-usb-libvirt-hotplug.rules
SUBSYSTEM=="usb",RUN+="/opt/usb-libvirt-hotplug/usb-libvirt-hotplug.sh win10"
